# files/init.d/S99alist
#!/bin/sh /etc/rc.common
START=99
start() {
    /mnt/sda1/alist/alist restart
}
stop() {
   /mnt/sda1/alist/alist stop
}